# supervisor/send-email-queue.ini
[supervisord]
nodaemon=true

[program:send-email-queue]
process_name=%(program_name)s_%(process_num)02d
command=php $PROJECT_DIR/artisan queue:work redis --queue=send-email-queue --tries=3 --timeout=60
autostart=true
autorestart=true
user=$USER
numprocs=1
redirect_stderr=true